# files/owrt-tgbot.cron
* * * * * /usr/bin/ucode -S /usr/lib/owrt-tgbot/alerts/runner.uc
0 0 * * * /usr/bin/ucode -S /usr/lib/owrt-tgbot/cron/midnight.uc
